# recipes-support/battlehill/files/battlehill.service
[Unit]
Description=Battlehill
After=local-fs.target daylited.service

[Service]
ExecStart=/usr/bin/battlehill

[Install]
WantedBy=multi-user.target
