# recipes-support/daylite/files/daylited.service
[Unit]
Description=Daylited
After=local-fs.target

[Service]
ExecStart=/usr/bin/daylited

[Install]
WantedBy=multi-user.target
